# build/supervisor.conf
[supervisord]
nodaemon = true

[program:symlink]
command=ln -sf /var/www/links/links.py /var/www/locust/links/

[program:locust]
command=locust -f /var/www/locust/locust_file.py -H %(ENV_WEBSITE_DOMAIN)s
